# dist/systemd/evo.service.d/hardware-audio-privileges.conf
# Systemd unit drop-in: filesystem-namespace grants for the
# hardware.audio-config plugin's privileged-action surface.
#
# The framework's reference `evo.service.example` sets
# `ProtectSystem=strict`, which mounts the host's filesystem
# read-only for the service except for paths explicitly listed
# in `ReadWritePaths`. The mpd-restart drop-in widens that to
# include `/etc/evo`. This drop-in adds the four further paths
# the hardware-audio plugin writes through its narrow sudoers
# grant.
#
# The privileged-action set the plugin invokes is declared in
# `dist/sudoers.d/evo-hardware-audio` (the operator-signed
# `EVO_HARDWARE_AUDIO_WRITE` + `EVO_HARDWARE_AUDIO_MODDER`
# Cmnd_Aliases). Sudo admits the elevation; the kernel still
# refuses the write if the path is read-only in the service's
# mount namespace, because sudo cannot escape the namespace's
# mount policy. The visible failure is
# `tee: <path>: Read-only file system` even though the
# interactive shell on the same host writes the file freely.
#
# Each path is listed individually with the `-` prefix per
# systemd's documented optional-path semantics: a missing
# file on this board class (Rockchip's `/boot/hw_intf.conf`
# absent on Pi; Pi's `/boot/firmware/*` absent on older Pi
# images that use `/boot/`) does NOT refuse unit start.
# Operators on board classes the plugin does not target see
# the path silently skipped.
#
# Per-path rationale:
#
#   - `/boot/firmware/config.txt` — Pi Bookworm / Trixie default.
#     `select_dac` / `clear_dac` write the dtoverlay managed
#     block here.
#   - `/boot/config.txt` — older Pi images. Same writes target
#     this path when `/boot/firmware/` is absent.
#   - `/boot/hw_intf.conf` — Rockchip-class boards. The
#     RockchipProvider equivalent of the Pi dtoverlay block.
#   - `/etc/modules-load.d` — directory grant so the Pi-only
#     `evo-i2c-dev.conf` drop-in can be created on first
#     `select_dac` against an I2S DAC catalogue row.
#   - `/boot/firmware/overlays` — directory grant for the DIY
#     DTBO modder workflow's `register_overlay` /
#     `remove_overlay` writes (operator-signed allowlist gate
#     applies BEFORE the sudo invocation; the namespace grant
#     is the kernel-level last layer).
#
# Vendor distributions that exclude the hardware.audio plugin
# (e.g. shipping vendor-locked hardware where DAC selection is
# not operator-tunable) can drop this conf; the framework's
# strict namespace stays in place.
#
# Apply: drop this file at
# `/etc/systemd/system/evo.service.d/`, run
# `systemctl daemon-reload && systemctl restart evo.service`.

[Service]
ReadWritePaths=-/boot/firmware/config.txt
ReadWritePaths=-/boot/config.txt
ReadWritePaths=-/boot/hw_intf.conf
ReadWritePaths=-/etc/modules-load.d
ReadWritePaths=-/boot/firmware/overlays
